# etc/udev/rules.d/61-hdparm.rules
ACTION=="add", KERNEL=="sd[a-z]*", ATTR{queue/rotational}=="1", \
  RUN+="/usr/bin/hdparm -B254 -S0 -W1 /dev/%k"
